# simple-card-auth.service
[Unit]
Description=SimpleCardAuth
Requires=pcscd.service

[Service]
ExecStart=/root/SimpleCardAuth/auth-loop.sh

[Install]
Also=pcscd.service
